# whatsapp-ai.service
[Unit]
Description=WhatsApp AI Chatbot
After=network.target

[Service]
Type=simple
# Ganti 'ubuntu' dengan user yang sesuai (ubuntu untuk Ubuntu, ec2-user untuk Amazon Linux)
User=ubuntu
# Ganti path sesuai lokasi aplikasi Anda
WorkingDirectory=/home/ubuntu/whatsapp-bot
Environment="NODE_ENV=production"
EnvironmentFile=/home/ubuntu/whatsapp-bot/.env
# Gunakan 'which node' untuk cek path node
ExecStart=/usr/bin/node /home/ubuntu/whatsapp-bot/app.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whatsapp-ai

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
